# .config/systemd/user/weather.timer
[Unit]
Description=Update weather.service

[Timer]
OnBootSec=3min
OnUnitActiveSec=20min

[Install]
WantedBy=timers.target
